# tests/perf/supervisord.conf
[supervisord]
logfile=./supervisord.log
errorlog=./supervisord.log

[supervisorctl]

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:webhook]
command=python worker.py
